# resources/rails/mongrel_cluster_init.erb
#!/bin/bash
#
# Copyright (c) 2007 Bradley Taylor, bradley@railsmachine.com
# Modified by Matt Patterson for dpkg-tools
# 
# mongrel_cluster       Startup script for Mongrel clusters.
#
# chkconfig: - 85 15
# description: mongrel_cluster manages multiple Mongrel processes for use \
#              behind a load balancer.
#              

CONF_DIR=<%= conf_dir_path %>
CONF_FILE=$CONF_DIR/mongrel_cluster.yml
PID_DIR=<%= pidfile_dir_path %>
USER=<%= username %>

RETVAL=0

# Gracefully exit if the controller is missing.
which mongrel_cluster_ctl >/dev/null || exit 0

# Go no further if config file is missing.
[ -f "$CONF_DIR/mongrel_cluster.yml" ] || exit 0

case "$1" in
    start)
      # Create pid directory
      mkdir -p $PID_DIR
      chown $USER:$USER $PID_DIR

      mongrel_rails cluster::start -C $CONF_FILE
      RETVAL=$?
  ;;
    stop)
      mongrel_rails cluster::stop -C $CONF_FILE
      RETVAL=$?
  ;;
    restart)
      mongrel_rails cluster::restart -C $CONF_FILE
      RETVAL=$?
  ;;
    status)
      mongrel_rails cluster::status -C $CONF_FILE
      RETVAL=$?
  ;;
    *)
      echo "Usage: mongrel_cluster {start|stop|restart|status}"
      exit 1
  ;;
esac      

exit $RETVAL
